# flask_boost/project/deploy/supervisor.conf
[program:#{project}]
command=/var/www/#{project}/venv/bin/gunicorn -c deploy/gunicorn.conf wsgi:app
directory=/var/www/#{project}
user=root
autostart=true
autorestart=true
environment = MODE="PRODUCTION"