# deploy/supervisor.conf
[program:proj]
command=/var/www/proj/venv/bin/gunicorn -c deploy/gunicorn.conf proj:app
directory=/var/www/proj
user=root
autostart=true
autorestart=true