# deploy/supervisor.conf
[program:tm2]
command=/var/www/project/venv/bin/gunicorn -c deploy/gunicorn.conf project:app
directory=/var/www/project
user=root
autostart=true
autorestart=true